# extras.pine
#!/bin/sh
## this file must be sourced to operate in the correct directory

## install packages from alpine to install init scripts...etc, then replace binaries
## with the upstream ones
apkc docker

## DOCKER
DOCKER_VERSION=`last_version docker/docker`
wget -q https://get.docker.com/builds/Linux/x86_64/docker-${DOCKER_VERSION}.tgz -O docker.tgz
tar xf docker.tgz
mv docker/docker* usr/bin/
rm -rf docker docker.tgz

## ETCD
cat << 'EOF' >usr/bin/etcdctl
#!/bin/sh

exec docker exec -e ETCDCTL_API=${ETCDCTL_API:-3} -it $(docker ps -f name=etcd -q) etcdctl $@
EOF
cat << 'EOF' >etc/init.d/etcd
#!/sbin/openrc-run

## Stub service to wait for the etcd container to start
SVCNAME="etcd container"

depend() {
    need net docker
}

start() {
    ebegin "Waiting for $SVCNAME"
    timeout -t 30 /bin/busybox sh -c ' \
    until [ "`docker inspect -f {{.State.Running}} $(docker ps -q -f name=etcd)`" == "true" ]; do
        sleep 0.5;
    done;'
    eend $? "$SVCNAME did not start"
}
EOF

## REBOOT MANAGER
## - The clearlock scripts dequeues the node from the list stored on etcd after boot
## - the reboot command is wrapped with a script to also support reboot locking
cat << 'EOF' >etc/init.d/clearlock
#!/sbin/openrc-run

description="Create targets for ostree deployments"

depend()
{
        need etcd
}

start() {
        ebegin "Clearning etcd reboot lock"
        if [ find /sbin/reboot -type f ] ; then
            exec /usr/bin/timeout -s2 -t900 /sbin/reboot lock clear
        fi
        eend $?
}

stop() {
        return
}
EOF
chmod +x etc/init.d/clearlock
cp ../reboot.pine sbin/reboot
chmod +x sbin/reboot

## SUP
## - the sup command is wrapped to use a default Supfile script located in etc/Supfile
## for common utilities like bootstrapping an etcd cluster
mkdir /go
export GOPATH=/go
go get -u github.com/pressly/sup/cmd/sup
mv /go/bin/sup usr/bin/sup.bin
cp ../Supfile etc/
cat << 'EOF' >usr/bin/sup
#!/bin/sh
SUP_FILE=/etc/Supfile
if $(echo "$@" | grep -qE "\-f\s+[^\s]+\s*") ; then
    exec /usr/bin/sup.bin -e SUP_FILE=$SUP_FILE $@
else
    exec /usr/bin/sup.bin -e SUP_FILE=$SUP_FILE -f $SUP_FILE $@
fi
EOF
chmod +x usr/bin/sup
