# dist/scripts/clearlock
#!/sbin/openrc-run

description="Delete the current host from the reboot queue"

depend()
{
        need networking
}

start() {
        ebegin "Clearing reboot lock"
        set -a
        . /etc/network-environment
        set +a
        if find /sbin/reboot -type f &>/dev/null; then
            exec /bin/busybox timeout -s2 -t900 /sbin/reboot lock clear
        fi
        eend $?
}

stop() {
        return
}
