# dist/scripts/vardirs
#!/sbin/openrc-run

description="Create targets for ostree deployments"

depend()
{
        need localmount
}

start() {
        ebegin "Setting dirs"
        mkdir -m 0755 -p /var/cache/rc /var/cache/apk /var/tmp /var/mnt /var/opt /var/srv /var/roothome /sysroot/ostree /sysroot/tmp
        ln -srnf /sysroot/home /var/home &>/dev/null
        eend $?
}

stop() {
        return
}